# device/energy_db_wd/energy_db_wd.conf
# start energy_db_wd service
#

start on runlevel [0123456]

stop on runlevel [!$RUNLEVEL]

respawn

exec /usr/local/bin/energy_db_wd /dev/ttyUSB0

post-stop script
#respawn delay
sleep 10 
end script
